# ci/ansible/roles/deploy/templates/certbot.service.j2
[Unit]
Description={{ description }} Service

[Service]
Type=oneshot
ExecStart="/usr/bin/certbot.sh"
